# .config/autostart/plumber.desktop
[Desktop Entry]
Type=Application
Name=Plumber
TryExec=plumber
Exec=plumber
